# kilnpi.service
[Unit]
Description=KilnPi Controller
After=network.target

[Service]
Type=simple
ExecStart=%h/.local/bin/uv run uvicorn backend.main:app --host 0.0.0.0 --port 8000
WorkingDirectory=%h/kilnpi
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
